# meta-ssd1306/recipes-currenttime/currenttime/files/currenttime.service
[Unit]
Description=Current time on SSD1306

[Service]
ExecStart=/usr/bin/currenttime
Restart=always

[Install]
WantedBy=multi-user.target