# meta-ssd1306/recipes-ssd1306/ssd1306/files/ssd1306.service
[Unit]
Description=Load SSD1306 module

[Service]
ExecStart=/etc/init.d/ssd1306_load.sh

[Install]
WantedBy=multi-user.target